# logtools/oem-varlog.cron
#!/bin/bash

/home/pi/emonpi/logtools/oem-varlog.manage
